# systemd.service
[Unit]
Description=jRAT Daemon
After=syslog.target

[Service]
ExecStart=/path/to

[Install]
WantedBy=multi-user.target
